# config/wol.service
[Unit]
Description=Enable Wake On Lan

[Service]
Type=oneshot
ExecStart = /sbin/ethtool --change eno1 wol g

[Install]
WantedBy=basic.target
